# stratum/hal/bin/bmv2/debian/postinst
# Copyright 2019-present Open Networking Foundation
# SPDX-License-Identifier: Apache-2.0
systemctl enable stratum_bmv2.service
