# stratum/hal/bin/bmv2/debian/postinst
# Copyright 2018 PLVision Open Networking Foundation
# Copyright 2018-present Open Networking Foundation
# SPDX-License-Identifier: Apache-2.0
systemctl enable stratum_bmv2.service
