# stratum/hal/bin/bmv2/debian/postinst
#!/bin/bash
# Copyright 2018 PLVision Open Networking Foundation
# Copyright 2018-present Open Networking Foundation
# SPDX-License-Identifier: Apache-2.0

# If systemd is installed,
if command -v systemctl &> /dev/null; then
  # Enable Stratum service
  systemctl enable stratum_bmv2.service
fi
