# stratum/hal/bin/barefoot/deb/postinst
# Working directory and logs
mkdir -p /var/run/stratum
mkdir -p /var/log/stratum

# Setting up huge page for bf_switchd
echo "vm.nr_hugepages = 128" >> /etc/sysctl.conf
sysctl -p /etc/sysctl.conf
mkdir /mnt/huge
mount -t hugetlbfs nodev /mnt/huge

# Enables Stratum service
systemctl enable stratum_bf.service
