# stratum/hal/bin/barefoot/deb/prerm
# Stop Stratum service
systemctl enable stratum_bf.service
