# stratum/hal/bin/barefoot/deb/prerm
# Stop and disable Stratum service
systemctl stop stratum_bf.service
systemctl disable stratum_bf.service
