# stratum/hal/bin/barefoot/deb/stratum_bf.service
[Unit]
Description=Stratum service
ConditionPathExists=/usr/bin/stratum_bf

[Service]
Type=simple
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=stratum_bf
WorkingDirectory=/var/run/stratum
LogsDirectory=stratum
RuntimeDirectory=stratum
ConfigurationDirectory=stratum
RuntimeDirectoryPreserve=yes
EnvironmentFile=/etc/stratum/stratum_bf.env
ExecStart=/bin/bash /usr/bin/stratum-entrypoint
Restart=on-failure

[Install]
WantedBy=multi-user.target
