# stratum/hal/bin/barefoot/deb/stratum_bf.service
[Unit]
Description=Stratum
ConditionPathExists=/usr/bin/stratum_bf
ConditionPathExists=/etc/stratum/stratum_bf.env

[Service]
Type=simple
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=stratum_bf

WorkingDirectory=/var/run/stratum
EnvironmentFile=/etc/stratum/stratum_bf.environment

ExecStartPre=/bin/mkdir -p /var/run/stratum

ExecStart=/usr/local/bin/stratum_bf -flagfile=$FLAG_FILE

ExecStop=/usr/bin/pkill stratum_bf

Restart=on-failure

[Install]
WantedBy=multi-user.target
